# supervisor/compasse-bulk-worker.conf
[program:compasse-bulk-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/compasse-backend/artisan queue:work redis --queue=bulk-uploads,broadcasts,emails --sleep=3 --tries=1 --timeout=3600 --memory=512
autostart=true
autorestart=true
numprocs=2
user=deploy
redirect_stderr=true
stdout_logfile=/var/www/compasse-backend/storage/logs/bulk-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
